# docker/api/php/supervisor/messenger.ini
[program:message]
command                 = sh -c "ulimit -n 10000; exec /usr/local/bin/php /var/www/coin24/api/bin/console messenger:consume async -vv --time-limit=3600 --memory-limit=256M"
process_name            = ProcessMessageSymfony
numprocs                = 1
autostart               = true
autorestart             = true
;user                    = root
stdout_logfile          = /var/log/supervisor/message.info.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /var/log/supervisor/message.error.log
stderr_logfile_maxbytes = 1MB
